# tools/aquarium-ml.service
[Unit]
Description=Unified Aquarium ML Service (PID + Water Change)
After=network.target mariadb.service mosquitto.service
Wants=mariadb.service mosquitto.service

[Service]
Type=simple
User=aquarium
Group=aquarium
WorkingDirectory=/opt/aquarium
ExecStart=/usr/bin/python3 /opt/aquarium/aquarium_ml_service.py --service
Restart=always
RestartSec=10

# Environment variables (update these for your setup)
Environment="MQTT_BROKER=localhost"
Environment="MQTT_PORT=1883"
Environment="MQTT_USER="
Environment="MQTT_PASSWORD="
Environment="MQTT_TOPIC_PREFIX=aquarium"
Environment="DB_HOST=localhost"
Environment="DB_PORT=3306"
Environment="DB_NAME=aquarium"
Environment="DB_USER=aquarium"
Environment="DB_PASSWORD=aquarium"

# Logging
StandardOutput=append:/var/log/aquarium-ml.log
StandardError=append:/var/log/aquarium-ml.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
